# xzm/003-settings/etc/profile.d/mc.sh
#!/bin/sh

# Ash does not work as the subshell yet:
alias mc="mc -u"
