# xzm/003-settings/etc/profile.d/variables.sh
#!/bin/sh

# Export variables:
export \
autostart=/etc/xdg/openbox/autostart \
scripts=/opt/scripts \
profile=/home/guest/.mozilla/firefox/c3pp43bg.default \
pprofile=/opt/scripts/guest/.mozilla/firefox/c3pp43bg.default \
json=/etc/opt/chrome/policies/managed/chrome.json \
chflags=/opt/google/chrome/chrome-flags.conf \
DISPLAY=:0
